# scripts/systemd/backend.service
[Unit]
Description=Character Booth Backend Service
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/Phone-Booth
ExecStart=/usr/bin/env bash backend/run_backend.sh
Restart=on-failure
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target


